# sentinela_receiver/systemd/sentinela-receiver-prod.service
[Unit]
Description=Sentinela Alarm Receiver — PROD (Sentinela_V18, port 10001)
After=network.target

[Service]
Type=simple
User=egarza
WorkingDirectory=/home/egarza/sentinela_receiver
ExecStart=/usr/bin/python3 /home/egarza/sentinela_receiver/receiver.py /home/egarza/sentinela_receiver/config_prod.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
